# config/clawbox-gateway.service
[Unit]
Description=ClawBox OpenClaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=clawbox
WorkingDirectory=/home/clawbox
ExecStartPre=/home/clawbox/clawbox/scripts/gateway-pre-start.sh
ExecStart=/home/clawbox/.npm-global/bin/openclaw gateway --allow-unconfigured --bind loopback
Restart=always
RestartSec=5
Environment=HOME=/home/clawbox
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=1024
Environment=OPENCLAW_VERBOSE=1
EnvironmentFile=-/home/clawbox/clawbox/data/network.env

[Install]
WantedBy=multi-user.target
