# config/clawbox-setup.service
[Unit]
Description=ClawBox Setup Wizard Web Server
After=network-online.target clawbox-ap.service
Wants=network-online.target

[Service]
Type=simple
User=clawbox
WorkingDirectory=/home/clawbox/clawbox
ExecStart=/usr/bin/node production-server.js
Restart=always
RestartSec=3
EnvironmentFile=-/home/clawbox/clawbox/data/network.env
EnvironmentFile=-/home/clawbox/clawbox/.env
Environment=NODE_ENV=production
Environment=PORT=80
Environment=HOSTNAME=0.0.0.0
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/clawbox/clawbox /home/clawbox/.bun /home/clawbox/.npm-global /home/clawbox/.npm /home/clawbox/.openclaw /tmp

[Install]
WantedBy=multi-user.target
